# services/linkedin-browser.service
[Unit]
Description=LinkedIn Browser (Chromium CDP)
After=network-online.target
Wants=network-online.target

[Service]
# Security: run as dedicated unprivileged user (not root)
User=linkedin-browser
Group=linkedin-browser
ExecStart=/opt/chromium-1208/chrome-linux/chrome --headless=new --no-sandbox --disable-gpu --disable-crash-reporter --remote-debugging-port=18801 --remote-debugging-address=127.0.0.1 --user-data-dir=/home/linkedin-browser/browser-data about:blank
Restart=always
RestartSec=5
KillMode=process
Environment=HOME=/home/linkedin-browser

[Install]
WantedBy=default.target
